# scripts/random-screens/random-screens.service
[Unit]
Description=Set random images for splash screens
After=home.mount

[Service]
Type=simple
ExecStart=/usr/share/remarkable/scripts/set-random-sleep.sh

[Install]
WantedBy=multi-user.target
